# agents/com.amir.macrig.plist.tpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
 <key>Label</key><string>com.amir.macrig</string>
 <key>ProgramArguments</key>
 <array><string>/Applications/MacRig.app/Contents/MacOS/MacRig</string></array>
 <key>EnvironmentVariables</key>
 <dict><key>MACRIG_DIR</key><string>__MACRIG_DIR__</string></dict>
 <key>RunAtLoad</key><true/>
 <key>KeepAlive</key><false/>
</dict></plist>
